# entrypoints/supervisord-control-queue.conf
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord.log

[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --sleep=1 --tries=1 --timeout=10800 --queue=default
autostart=true
autorestart=true
numprocs=%(ENV_NUMPROCS_QUEUE)s
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:queue-export]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --sleep=1 --tries=1 --timeout=10800 --queue=export
autostart=true
autorestart=true
numprocs=%(ENV_NUMPROCS_QUEUE_EXPORT)s
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
